# systemd/antiddos-ssh.service
[Unit]
Description=Anti-DDoS SSH Protection Service
After=network.target sshd.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 -m antiddos.ssh_protection -c /etc/antiddos/config.yaml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=false
ProtectHome=true

[Install]
WantedBy=multi-user.target
